# roles/apt-mirror/templates/units/update-apt-mirror.service.j2
[Unit]
Description=apt-mirror update service

[Service]
Type=oneshot
WorkingDirectory={{ apt_mirror_path  }}
ExecStart=/usr/bin/apt-mirror
